# deploy/kylin/systemd/loong64-b1-upgrade.service
[Unit]
Description=loong64-b1-go system upgrade
Documentation=https://github.com/kenichiLyon/loong64-b1-go
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=loong64b1
Group=loong64b1
WorkingDirectory=/opt/loong64-b1-go
EnvironmentFile=/etc/loong64-b1-go/loong64-b1-go.env
ExecStart=/opt/loong64-b1-go/bin/loong64-b1-upgrade-linux-loong64 up
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/loong64-b1-go /var/log/loong64-b1-go
